# usr/share/keepalive/docs/debian-keepalive/prerm
# deb prerm 2017-01-26
thisservice=keepalive.service
{
case "${1}" in
   remove)
      # Final removal.
      systemctl stop "${thisservice}"
      systemctl disable "${thisservice}"
      ;;
   upgrade)
      # Upgrade. $2 is new-version-number
      :
      ;;
esac
} 1>/dev/null 2>&1
exit 0 
